# config/fujiapp.service
[Service]
WorkingDirectory=/usr/src/app
ExecStart=/usr/bin/npm start
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=fujiapp
User=nobody
Group=nogroup
Environment='NODE_ENV=production'

[Install]
WantedBy=multi-user.target
